# host/systemd/gopro-autorevive.service.in
[Unit]
Description=GoPro auto-revive watcher (power-cycle a camera that fell off the USB bus)
After=docker.service
Wants=docker.service

[Service]
Type=simple
SyslogIdentifier=revive
User=__USER__
# Paths filled in by install.sh (sed on __GOPRO_SCRIPTS__) -- never hard-coded.
ExecStart=/usr/bin/stdbuf -oL __GOPRO_SCRIPTS__/revive.sh --watch
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
