# host/systemd/gopro-manager.service.in
[Unit]
Description=GoPro manager (arm cameras, recording, watchdog/EMERGENCY signal)
After=docker.service network-online.target
Wants=docker.service

[Service]
# manager_up.sh launches a DETACHED docker container and returns, so this is a
# oneshot that stays "active" afterwards; the container itself (and its docker
# --restart unless-stopped) keeps it alive at runtime, systemd (re)starts it at boot.
Type=oneshot
RemainAfterExit=yes
User=__USER__
Environment=HOME=__HOME__
# Paths filled in by install.sh (sed on __GOPRO_SCRIPTS__) -- never hard-coded.
ExecStart=__GOPRO_SCRIPTS__/manager_up.sh
ExecStop=__GOPRO_SCRIPTS__/manager_down.sh

[Install]
WantedBy=multi-user.target
